# saltstack/salt/dropbox/files/etc/systemd/system/dropbox@.service.d/override.conf
[Service]
Environment=DISPLAY=:0
SuccessExitStatus=KILL
